# autolight/doc/systemd/autolight.service

[Unit]
Description=Start Autolight LED security light
After=ssh.service

[Service]
EnvironmentFile=/home/al/autolight/environment.file
ExecStart=/home/al/autolight/securitylight.py
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=autolight
User=al
Group=al

[Install]
WantedBy=multi-user.target
